# BLIP/services/BLIP.service
[Unit]
Description=BLIP Discord Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=0
User=ubuntu
WorkingDirectory=/home/ubuntu/BLIP
ExecStart=/home/ubuntu/BLIP/discord.sh

[Install]
WantedBy=default.target
