# moondream/services/moondream-api.service
[Unit]
Description=Moondream REST API Service
After=network.target
StartLimitBurst=3
StartLimitIntervalSec=300

[Service]
Type=simple
Restart=on-failure
RestartSec=5
User=sd
Group=sd
WorkingDirectory=/home/sd/animal-farm/moondream
ExecStart=/home/sd/animal-farm/moondream/start.sh
EnvironmentFile=/home/sd/animal-farm/moondream/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
