# ocr/services/ocr-api.service
[Unit]
Description=PaddleOCR API Service (Modern OCR replacement)
After=network.target
Conflicts=tesseract-api.service

[Service]
Type=simple
User=sd
Group=sd
WorkingDirectory=/home/sd/ocr
ExecStart=/home/sd/ocr/venv/bin/python REST.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
